# battery-monitor/battery-monitor.timer
[Timer]
OnUnitInactiveSec=10

[Install]
WantedBy=timers.target
